# profile/profile.d/1password.sh
alias oplogin='op signin ${ONE_PASSWORD_DOMAIN}'
