# profile/profile.d/python.sh
eval "$("${BREW_PREFIX}/bin/pyenv" init -)"

export PATH="${HOME}/.local/bin:${PATH}"
