# profile/profile.d/ruby.sh
_pathmunge "${BREW_PREFIX}/lib/ruby/gems/3.3.0/bin"
_pathmunge "${BREW_PREFIX}/opt/ruby/bin"

