# linux/profile/profile.d/go.sh
#!/usr/bin/env bash

export GOPATH="${HOME}/source/go"
export PATH="${GOPATH}/bin:${PATH}"

