# linux/profile/profile.d/profile.d
/volumes/user/source/.dotfiles/linux/profile/profile.d